# RPiLed/ledTelegramBot-watcher.path
[Path]
PathModified=/home/pi/Led

[Install]
WantedBy=multi-user.target